# packaging/systemd/arcwell-worker.service.in
[Unit]
Description=Arcwell local worker
After=network-online.target

[Service]
Type=simple
Environment="@ARCWELL_HOME_ENV@"
WorkingDirectory="@ARCWELL_HOME_PATH@"
ExecStart="@ARCWELL_BIN_EXEC@" worker run --max-jobs-per-tick 10 --idle-sleep-ms 5000
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths="@ARCWELL_HOME_PATH@"
StandardOutput=append:@LOG_DIR_PATH@/worker.log
StandardError=append:@LOG_DIR_PATH@/worker.err.log

[Install]
WantedBy=default.target
